# systemd/video-bot.timer
[Unit]
Description=Run Video Bot Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=video-bot.service

[Install]
WantedBy=timers.target
